# etc/hop3-nginx.path
# Goes to: /etc/systemd/system/hop3-nginx.path

[Unit]
Description=Monitor ~/nginx for changes

[Path]
PathChanged=/home/hop3/nginx
Unit=hop3-nginx.service

[Install]
WantedBy=multi-user.target
